# .config/supervisord.conf
[program:nginx]
command=nginx -g 'daemon off;'

[program:uwsgi]
command=uwsgi --ini /srv/project/.config/uwsgi.ini

